# systemd/tekton-data-check.service
[Unit]
Description=Tekton Market Data Freshness Check
After=network.target

[Service]
Type=oneshot
User=tony
WorkingDirectory=/home/tony/tekton-ai-trader
EnvironmentFile=/home/tony/tekton-ai-trader/.env
ExecStart=/home/tony/tekton-ai-trader/venv/bin/python /home/tony/tekton-ai-trader/check_data_freshness.py
StandardOutput=journal
StandardError=journal
